# scripts/runit/mcp-server/run
#!/data/data/com.termux/files/usr/bin/sh
# runit service for Termux MCP Server (Rust)
set -eu
exec 2>&1

TOKEN_FILE="$HOME/.termux_mcp_token"

if [ ! -s "$TOKEN_FILE" ]; then
  echo "ERROR: $TOKEN_FILE is missing or empty. Create it with a strong bearer token before starting the service." >&2
  exit 1
fi

TOKEN="$(tr -d '\r\n' < "$TOKEN_FILE")"
if [ -z "$(printf '%s' "$TOKEN" | tr -d '[:space:]')" ]; then
  echo "ERROR: $TOKEN_FILE contains only whitespace. Replace it with a strong bearer token." >&2
  exit 1
fi

export MCP__SERVER__HOST="127.0.0.1"
export MCP__SERVER__PORT="8000"
export MCP__AUTH__STATIC_TOKEN="$TOKEN"

exec "$HOME/bin/termux-mcp-server"
